# bixi-display.service
[Unit]
Description=BIXI Pixoo 64 display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/bixi-display
Environment=BIXI_CONFIG_PATH=/opt/bixi-display/config.json
Environment=BIXI_LOG_PATH=/var/log/bixi-display.log
ExecStart=/opt/bixi-display/.venv/bin/python /opt/bixi-display/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
